# .config/systemd/user/lemonade.service
[Unit]
Description=Lemonade

[Service]
ExecStart=/home/seth/.local/bin/lemonade server

[Install]
WantedBy=default.target
